# docker/hbase/config/supervisor.conf
[program:master]
command=/bin/sh -c "/opt/hbase/bin/hbase master start"
stdout_logfile = /data/logs/hbase-master.log
redirect_stderr=true
autostart=true
autorestart=false
priority=10
